# mydockerbuild/ssh/debian/supervisor/nodejs/supervisor-config/node.sv.conf
[program:node]
command = /usr/bin/node /home/udemy/app.js
